# roles/gh-runner/templates/gh-runner.service.j2
[Unit]
Description=GitHub action self hosted runner

[Service]
User=gh-runner
TimeoutStartSec=0
Restart=always
ExecStart=/home/gh-runner/{{ item.name }}/run.sh

[Install]
WantedBy=default.target
